# cadocms/templates/config/crontab
# This user crontab file was generated by CadoCMS.
# Please don't edit it manually
# Instead please copy "crontab" file 
{% for SITE in SITES %}
# from cadocms/teamplates/config/ to {{SITE.CADO_PROJECT}}/templates/config/
{% endfor %}
# and add anything you want.
# 
# This file has normal crontab syntax.
# Example entry:
# 0 5 * * 1 ls -l /home/
# 
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command

{% for SITE in SITES %}



@reboot cd {{HOST.SRCROOT}}; {{HOST.PYTHON_PREFIX}}python manage.py runfcgi {{SITE.CADO_PROJECT}} maxchildren={% if HOST.CLASS = 'TEST' %}1{% else %}3{% endif %} method=threaded socket={{HOST.APPROOT}}{{SITE.CADO_PROJECT}}.sock pidfile={{HOST.APPROOT}}{{SITE.CADO_PROJECT}}.pid; sleep 5; chmod a+rwx {{HOST.APPROOT}}{{SITE.CADO_PROJECT}}.sock
{% endfor %}

3 0 * * * kill -9 `cat {{HOST.APPROOT}}{{SITE.CADO_PROJECT}}.pid`; cd {{HOST.SRCROOT}}; {{HOST.PYTHON_PREFIX}}python manage.py runfcgi {{SITE.CADO_PROJECT}} maxchildren={% if HOST.CLASS = 'TEST' %}1{% else %}3{% endif %} method=threaded socket={{HOST.APPROOT}}{{SITE.CADO_PROJECT}}.sock pidfile={{HOST.APPROOT}}{{SITE.CADO_PROJECT}}.pid; sleep 5; chmod a+rwx {{HOST.APPROOT}}{{SITE.CADO_PROJECT}}.sock
{% endfor %}
